# php-worker/supervisord.d/laravel-worker-restart.conf
[program:periodicrestarter]
command = sh -c "php /var/www/laravel/artisan queue:restart;sleep 600"
autorestart = true